# init/couchpotato.service
[Unit]
Description=CouchPotato application instance
After=network.target

[Service]
ExecStart=/var/lib/CouchPotatoServer/CouchPotato.py
Type=simple
User=couchpotato
Group=couchpotato
Restart=always
RestartSec=2s

[Install]
WantedBy=multi-user.target
